# scripts/systemd/algotrader-dashboard.service
[Unit]
Description=AlgoTrader Dashboard
After=network.target postgresql.service algotrader.service
Wants=algotrader.service

[Service]
Type=simple
User=algotrader
Group=algotrader
WorkingDirectory=/opt/algotrader
Environment=PYTHONPATH=/opt/algotrader
Environment=PATH=/opt/algotrader/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=DATABASE_URL=postgresql://localhost/algotrader

ExecStartPre=/opt/algotrader/.venv/bin/python -c "import algotrader"
ExecStart=/opt/algotrader/.venv/bin/python -m algotrader.cli dashboard
ExecReload=/bin/kill -HUP $MAINPID

Restart=on-failure
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/algotrader/data /opt/algotrader/logs /opt/algotrader/output
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
